# manjaro/gnome-next/live-overlay/etc/profile.d/gsk_renderer.sh
# Test for live ISO to avoid failure with KVM

export GSK_RENDERER=ngl
